# supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:llama-server]
command=/usr/local/bin/llama-server
    -m "%(ENV_LLAMA_MODEL_PATH)s"
    --mmproj "%(ENV_LLAMA_MMPROJ_PATH)s"
    -c %(ENV_LLAMA_CONTEXT_LENGTH)s
    --host %(ENV_LLAMA_SERVER_HOST)s
    --port %(ENV_LLAMA_SERVER_PORT)s
    --reasoning-budget 0
    --reasoning off
    -ngl %(ENV_LLAMA_NGL)s
    -np 1
priority=1
autostart=true
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:fastapi]
command=uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
directory=/app/backend
priority=2
autostart=true
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
priority=3
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
